# claude-wrapper.timer
# Optional systemd timer for periodic cleanup
# Runs cleanup daily at 2 AM

[Unit]
Description=Claude Wrapper Daily Cleanup Timer
Documentation=https://github.com/yourusername/claude-wrapper
Requires=claude-wrapper.service

[Timer]
OnCalendar=daily
OnCalendar=02:00
Persistent=true

[Install]
WantedBy=timers.target
